# tools/mempalace/com.user.mempalace-snapshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS launchd equivalent of mempalace-snapshot.timer + .service.
     Snapshots the transactional chroma.sqlite3 every 6h. __HOME__ is rendered
     to $HOME by init.sh. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.mempalace-snapshot</string>
  <key>ProgramArguments</key>
  <array>
    <string>__HOME__/.local/bin/mempalace-snapshot.sh</string>
  </array>
  <!-- Every 6 hours (21600s). RunAtLoad gives a catch-up shortly after login. -->
  <key>StartInterval</key>
  <integer>21600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__HOME__/.mempalace/logs/snapshot.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.mempalace/logs/snapshot.log</string>
  <key>Nice</key>
  <integer>10</integer>
</dict>
</plist>
