# tools/mempalace/mempalace-prune.timer
[Unit]
Description=Daily mempalace junk-drawer prune

[Timer]
# Daily 03:47 local; off the :00 mark. Persistent catches up if the box was off.
# Daily (not weekly) keeps the recall pollution window short: the Stop hook
# re-mines the whole session dir in convos mode (which ignores .gitignore and
# has no exclude), so junk can only be removed after ingest — keep that lag <1d.
OnCalendar=*-*-* 03:47:00
Persistent=true

[Install]
WantedBy=timers.target
