# tools/graphify/com.user.graphify-reseed.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS launchd equivalent of graphify-reseed.timer + .service.
     Nightly graphify→mempalace structural reseed. __HOME__ is rendered to $HOME,
     and __GRAPHIFY_REPO_ARGS__ is replaced with one <string> entry per repo in
     init.sh's GRAPHIFY_REPOS list, by init.sh at install time. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.graphify-reseed</string>
  <key>ProgramArguments</key>
  <array>
    <string>__HOME__/.local/bin/graphify-reseed.sh</string>
__GRAPHIFY_REPO_ARGS__
  </array>
  <!-- No Weekday key => fires daily at the given Hour/Minute. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>4</integer>
    <key>Minute</key><integer>13</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__HOME__/.mempalace/logs/graphify-reseed.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.mempalace/logs/graphify-reseed.log</string>
  <key>Nice</key>
  <integer>10</integer>
</dict>
</plist>
