# tools/graphify/graphify-reseed.service
[Unit]
Description=graphify→mempalace structural reseed (wipe-and-replace)
Documentation=https://github.com/safishamsi/graphify

[Service]
Type=oneshot
# __GRAPHIFY_REPOS__ is rendered to the space-separated GRAPHIFY_REPOS list by init.sh.
ExecStart=%h/.local/bin/graphify-reseed.sh __GRAPHIFY_REPOS__
Nice=10
