# tools/graphify/graphify-reseed.service
[Unit]
Description=graphify→mempalace structural reseed (wipe-and-replace)
Documentation=https://github.com/safishamsi/graphify

[Service]
Type=oneshot
# __REPO__ is rendered to the repo path by init.sh at install time.
ExecStart=%h/.local/bin/graphify-reseed.sh __REPO__
Nice=10
