# tools/graphify/graphify-reseed.timer
[Unit]
Description=Nightly graphify→mempalace structural reseed

[Timer]
# Nightly 04:13 local; off the :00 mark and after the 03:47 prune. Persistent
# catches up if the box was off. Runs out of session: a live Claude session holds
# the mempalace write-lock, so the reseed skips (MCP guard) and retries next night.
OnCalendar=*-*-* 04:13:00
Persistent=true

[Install]
WantedBy=timers.target
