# tools/graphify/graphify-reseed.service
[Unit]
Description=graphifyŌåÆmempalace structural reseed (wipe-and-replace)
Documentation=https://github.com/safishamsi/graphify

[Service]
Type=oneshot
# systemd user services also start with a minimal PATH; add ~/.local/bin for graphify/mempalace.
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
# __GRAPHIFY_REPOS__ is rendered to the space-separated GRAPHIFY_REPOS list by init.sh.
ExecStart=%h/.local/bin/graphify-reseed.sh __GRAPHIFY_REPOS__
Nice=10
