# tools/chatgpt-toggle/systemd/chatgpt-bridge.service
[Unit]
Description=claude-code-proxy (ChatGPT/Codex bridge for the Claude CLI GPT toggle)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=PORT=18765
# Force reasoning effort for all Codex/ChatGPT requests (low|medium|xhigh|max).
Environment=CCP_CODEX_EFFORT=medium
ExecStart=%h/.local/bin/claude-code-proxy serve --no-monitor
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
