# tools/chatgpt-toggle/systemd/chatgpt-model-refresh.timer
[Unit]
Description=Periodically refresh the chatgpt-toggle default model (auto-pick new ChatGPT models)

[Timer]
OnBootSec=2min
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
